# rootfs/etc/s6-overlay/s6-rc.d/ocserv/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
s6-setuidgid ocserv
/usr/sbin/ocserv --foreground --config /etc/ocserv/ocserv.conf --log-stderr

