# rootfs/etc/s6-overlay/s6-rc.d/00-init-users/up
#!/command/execlineb -P
with-contenv
/etc/s6-overlay/s6-rc.d/00-init-users/run
