# rootfs/etc/s6-overlay/s6-rc.d/init-config/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

echo "[INIT-CONFIG] Ensuring ocserv configuration file exists"

# Create default config if missing
if [ ! -f /etc/ocserv/ocserv.conf ]; then
    cp /usr/share/ocserv/ocserv.conf.template /etc/ocserv/ocserv.conf
    chown "$PUID:$PGID" /etc/ocserv/ocserv.conf
    echo "[INIT-CONFIG] Copied default config to /etc/ocserv/ocserv.conf"
fi

echo "[INIT-CONFIG] ocserv configuration file is present"
