# rootfs/etc/s6-overlay/s6-rc.d/svc-ocserv/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

log "[SVC-OCSERV] Starting ocserv service"

# Ensure runtime directories exist
mkdir -p /run/ocserv /var/log/ocserv /var/lib/ocserv || true
chmod 0755 /run/ocserv /var/log/ocserv /var/lib/ocserv || true

exec fdmove -c 2 1 /usr/sbin/ocserv --foreground --config /etc/ocserv/ocserv.conf --log-stderr
