# rootfs/etc/s6-overlay/s6-rc.d/init-nat/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

log "[INIT-NAT] Initializing NAT and forwarding settings"

# Enable forwarding - check first, then try to set if needed
ipv4_forward=$(cat /proc/sys/net/ipv4/ip_forward)
if [ "$ipv4_forward" != "1" ]; then
    if echo 1 > /proc/sys/net/ipv4/ip_forward 2>/dev/null; then
        log "[INIT-NAT] IPv4 forwarding enabled"
    else
        log_error "[INIT-NAT] Warning: IPv4 forwarding disabled. Use --sysctl net.ipv4.ip_forward=1"
        sleep infinity
    fi
fi

if [ "$ipv6_forward" = "1" ]; then
    ipv6_forward=$(cat /proc/sys/net/ipv6/conf/all/forwarding)
    if [ "$ipv6_forward" != "1" ]; then
        if echo 1 > /proc/sys/net/ipv6/conf/all/forwarding 2>/dev/null; then
            log "[INIT-NAT] IPv6 forwarding enabled"
        else
            log_error "[INIT-NAT] Warning: IPv6 forwarding disabled. Use --sysctl net.ipv6.conf.all.forwarding=1"
        fi
    fi
fi

# Translate the iptables-style interface wildcard (e.g. vpns+) to nftables syntax (vpns*)
case "$vpn_if" in
    *+) vpn_ifname="${vpn_if%+}*" ;;
    *)  vpn_ifname="$vpn_if" ;;
esac

# Gateway egress interfaces: an upstream gateway (VPN_GATEWAY / VPN_GATEWAYS)
# may sit on a different interface than the WAN - e.g. a bridge to a VPN
# sidecar container while the WAN is a macvlan ISP uplink. Masquerade and
# forward via each such interface too, so the upstream sees this container's
# address and needs no return route to the client subnet.
gw_dev() { # print the egress interface for a gateway IP ($1 = -4|-6, $2 = IP)
    ip "$1" route get "$2" 2>/dev/null \
        | awk '{for (i = 1; i < NF; i++) if ($i == "dev") {print $(i + 1); exit}}'
}

gw_ifs4=""
for gw in $vpn_gateway $(kv_pairs "$vpn_gateways" | awk '{print $2}'); do
    dev="$(gw_dev -4 "$gw")"
    if [ -z "$dev" ]; then
        log_error "[INIT-NAT] Warning: cannot resolve egress interface for gateway $gw"
    elif [ "$dev" != "$wan_if" ]; then
        case " $gw_ifs4 " in
            *" $dev "*) ;;
            *) gw_ifs4="$gw_ifs4 $dev" ;;
        esac
    fi
done

gw_ifs6=""
for gw in $vpn_gateway6 $(kv_pairs "$vpn_gateways6" | awk '{print $2}'); do
    dev="$(gw_dev -6 "$gw")"
    if [ -z "$dev" ]; then
        log_error "[INIT-NAT] Warning: cannot resolve egress interface for gateway $gw"
    elif [ "$dev" != "$wan_if" ]; then
        case " $gw_ifs6 " in
            *" $dev "*) ;;
            *) gw_ifs6="$gw_ifs6 $dev" ;;
        esac
    fi
done

# Masqueraded source subnets: the VPN client subnet plus any forwarded LAN
# subnets (FORWARD_SUBNETS / FORWARD_SUBNETS6) that use this container as
# their gateway. Each subnet is masqueraded via the WAN and via every
# gateway egress interface.
masq_subnets4="$vpn_subnet"
for s in $(fwd_subnet_pairs "$forward_subnets" | awk '{print $1}'); do
    masq_subnets4="$masq_subnets4 $s"
done
masq_subnets6="$ipv6_subnet"
for s in $(fwd_subnet_pairs "$forward_subnets6" | awk '{print $1}'); do
    masq_subnets6="$masq_subnets6 $s"
done

nat_rules=""
for dev in $wan_if $gw_ifs4; do
    for s in $masq_subnets4; do
        log "[INIT-NAT] Masquerading $s via $dev"
        nat_rules="$nat_rules
        ip saddr $s oifname \"$dev\" masquerade"
    done
done
if [ "$ipv6_nat" = "1" ]; then
    for dev in $wan_if $gw_ifs6; do
        for s in $masq_subnets6; do
            log "[INIT-NAT] Masquerading $s via $dev (IPv6)"
            nat_rules="$nat_rules
        ip6 saddr $s oifname \"$dev\" masquerade"
        done
    done
fi

gw_fwd_rules=""
gw_fwd_seen=""
for dev in $gw_ifs4 $gw_ifs6; do
    case " $gw_fwd_seen " in
        *" $dev "*) continue ;;
    esac
    gw_fwd_seen="$gw_fwd_seen $dev"
    gw_fwd_rules="$gw_fwd_rules
        iifname \"$vpn_ifname\" oifname \"$dev\" accept
        iifname \"$dev\" oifname \"$vpn_ifname\" ct state related,established accept"
done

# Apply the ruleset in a dedicated table. Deleting it first keeps this idempotent
# across service restarts. Forwarding rules in the inet family cover IPv4 and IPv6.
log "[INIT-NAT] Setting up forwarding rules"
nft delete table inet ocserv 2>/dev/null || true
nft -f - <<EOF
table inet ocserv {
    chain forward {
        type filter hook forward priority 0; policy accept;
        iifname "$vpn_ifname" oifname "$wan_if" accept
        iifname "$wan_if" oifname "$vpn_ifname" ct state related,established accept$gw_fwd_rules
    }
    chain postrouting {
        type nat hook postrouting priority 100; policy accept;$nat_rules
    }
}
EOF

log "[INIT-NAT] NAT and forwarding setup complete"
