# rootfs/etc/s6-overlay/s6-rc.d/init-nat/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

echo "[INIT-NAT] Initializing NAT and forwarding settings"

# Enable forwarding - check first, then try to set if needed
ipv4_forward=$(cat /proc/sys/net/ipv4/ip_forward)
if [ "$ipv4_forward" != "1" ]; then
    if echo 1 > /proc/sys/net/ipv4/ip_forward 2>/dev/null; then
        echo "[INIT-NAT] IPv4 forwarding enabled"
    else
        echo "[INIT-NAT] Warning: IPv4 forwarding disabled. Use --sysctl net.ipv4.ip_forward=1" >&2
        sleep infinity
    fi
fi

if [ "$IPV6_FORWARD" = "1" ]; then
    ipv6_forward=$(cat /proc/sys/net/ipv6/conf/all/forwarding)
    if [ "$ipv6_forward" != "1" ]; then
        if echo 1 > /proc/sys/net/ipv6/conf/all/forwarding 2>/dev/null; then
            echo "[INIT-NAT] IPv6 forwarding enabled"
        else
            echo "[INIT-NAT] Warning: IPv6 forwarding disabled. Use --sysctl net.ipv6.conf.all.forwarding=1" >&2
        fi
    fi
fi

# IPv4 NAT
echo "[INIT-NAT] Setting up IPv4 NAT for subnet $VPN_SUBNET via $WAN_IF"
iptables -t nat -C POSTROUTING -s "$VPN_SUBNET" -o "$WAN_IF" -j MASQUERADE 2>/dev/null || \
    iptables -t nat -A POSTROUTING -s "$VPN_SUBNET" -o "$WAN_IF" -j MASQUERADE

# IPv4 forwarding rules
echo "[INIT-NAT] Setting up IPv4 forwarding rules"
iptables -C FORWARD -i "$VPN_IF" -o "$WAN_IF" -j ACCEPT 2>/dev/null || \
    iptables -A FORWARD -i "$VPN_IF" -o "$WAN_IF" -j ACCEPT
iptables -C FORWARD -i "$WAN_IF" -o "$VPN_IF" -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT 2>/dev/null || \
    iptables -A FORWARD -i "$WAN_IF" -o "$VPN_IF" -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT

# IPv6 NAT (optional)
if [ "$IPV6_NAT" = "1" ]; then
    echo "[INIT-NAT] Setting up IPv6 NAT for subnet $IPV6_SUBNET via $WAN_IF"
    ip6tables -t nat -C POSTROUTING -s "$IPV6_SUBNET" -o "$WAN_IF" -j MASQUERADE 2>/dev/null || \
        ip6tables -t nat -A POSTROUTING -s "$IPV6_SUBNET" -o "$WAN_IF" -j MASQUERADE
    echo "[INIT-NAT] Setting up IPv6 forwarding rules"
    ip6tables -C FORWARD -i "$VPN_IF" -o "$WAN_IF" -j ACCEPT 2>/dev/null || \
        ip6tables -A FORWARD -i "$VPN_IF" -o "$WAN_IF" -j ACCEPT
    ip6tables -C FORWARD -i "$WAN_IF" -o "$VPN_IF" -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT 2>/dev/null || \
        ip6tables -A FORWARD -i "$WAN_IF" -o "$VPN_IF" -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT
fi

echo "[INIT-NAT] NAT and forwarding setup complete"
