# rootfs/etc/s6-overlay/s6-rc.d/init-nat/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eu

. /usr/local/bin/backend-functions

log "[INIT-NAT] Initializing NAT and forwarding settings"

# Enable forwarding - check first, then try to set if needed
ipv4_forward=$(cat /proc/sys/net/ipv4/ip_forward)
if [ "$ipv4_forward" != "1" ]; then
    if echo 1 > /proc/sys/net/ipv4/ip_forward 2>/dev/null; then
        log "[INIT-NAT] IPv4 forwarding enabled"
    else
        log_error "[INIT-NAT] Warning: IPv4 forwarding disabled. Use --sysctl net.ipv4.ip_forward=1"
        sleep infinity
    fi
fi

if [ "$ipv6_forward" = "1" ]; then
    ipv6_forward=$(cat /proc/sys/net/ipv6/conf/all/forwarding)
    if [ "$ipv6_forward" != "1" ]; then
        if echo 1 > /proc/sys/net/ipv6/conf/all/forwarding 2>/dev/null; then
            log "[INIT-NAT] IPv6 forwarding enabled"
        else
            log_error "[INIT-NAT] Warning: IPv6 forwarding disabled. Use --sysctl net.ipv6.conf.all.forwarding=1"
        fi
    fi
fi

# Translate the iptables-style interface wildcard (e.g. vpns+) to nftables syntax (vpns*)
case "$vpn_if" in
    *+) vpn_ifname="${vpn_if%+}*" ;;
    *)  vpn_ifname="$vpn_if" ;;
esac

# IPv4 NAT
log "[INIT-NAT] Setting up IPv4 NAT for subnet $vpn_subnet via $wan_if"
nat_rules="ip saddr $vpn_subnet oifname \"$wan_if\" masquerade"

# IPv6 NAT (optional)
if [ "$ipv6_nat" = "1" ]; then
    log "[INIT-NAT] Setting up IPv6 NAT for subnet $ipv6_subnet via $wan_if"
    nat_rules="$nat_rules
        ip6 saddr $ipv6_subnet oifname \"$wan_if\" masquerade"
fi

# Apply the ruleset in a dedicated table. Deleting it first keeps this idempotent
# across service restarts. Forwarding rules in the inet family cover IPv4 and IPv6.
log "[INIT-NAT] Setting up forwarding rules"
nft delete table inet ocserv 2>/dev/null || true
nft -f - <<EOF
table inet ocserv {
    chain forward {
        type filter hook forward priority 0; policy accept;
        iifname "$vpn_ifname" oifname "$wan_if" accept
        iifname "$wan_if" oifname "$vpn_ifname" ct state related,established accept
    }
    chain postrouting {
        type nat hook postrouting priority 100; policy accept;
        $nat_rules
    }
}
EOF

log "[INIT-NAT] NAT and forwarding setup complete"
